# backup/crontab
0 23 1 * * /usr/local/bin/run-backup.sh >> /proc/1/fd/1 2>&1
